# deploy/callmor-api.service
[Unit]
Description=Callmor API Server (Axum)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=pasharobert
WorkingDirectory=/home/pasharobert/callmor-remote-desktop
EnvironmentFile=/home/pasharobert/callmor-remote-desktop/.env
ExecStart=/home/pasharobert/callmor-remote-desktop/target/release/callmor-api
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=callmor-api

[Install]
WantedBy=multi-user.target
