# packaging/deb/postinst
#!/bin/bash
set -e

# Create config directory if it doesn't exist
mkdir -p /etc/callmor-agent

# Only install default config if none exists (don't overwrite a real config on upgrade)
if [ ! -f /etc/callmor-agent/agent.conf ]; then
    cp /usr/share/callmor-agent/agent.conf.template /etc/callmor-agent/agent.conf
    chmod 600 /etc/callmor-agent/agent.conf
fi

# Reload systemd and start/enable the service. Agent will self-enroll on first
# run using the token baked into agent.conf by the downloader.
systemctl daemon-reload
systemctl enable callmor-agent >/dev/null 2>&1 || true
systemctl restart callmor-agent >/dev/null 2>&1 || true

echo ""
echo "================================================"
echo " Callmor Agent installed and started."
echo " Status: sudo systemctl status callmor-agent"
echo " Logs:   sudo journalctl -u callmor-agent -f"
echo "================================================"
echo ""

exit 0
